# units/pid-failsafe.service
[Unit]
Description=PID failsafe: re-open SSH if Tailscale not ready
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/pid-failsafe.sh
