# units/pid-failsafe.timer
[Unit]
Description=Run PID failsafe 10 min after boot

[Timer]
OnBootSec=10min
Persistent=true

[Install]
WantedBy=timers.target
