# units/pid-selftest.service
[Unit]
Description=PID node self-test
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/pid-selftest.sh
StandardOutput=journal
StandardError=journal
